# user/nsm-generic-mixer-session-save.timer
[Unit]
Description=Timer to launch NSM save

[Timer]
OnActiveSec=2min
OnUnitActiveSec=2min

[Install]
WantedBy=timers.target

