# user/wallpaper-restore.service
[Unit]
Description=Restore wallpaper

[Service]
Type=oneshot
ExecStart=/usr/bin/nitrogen --restore

[Install]
WantedBy=graphical-session.target
